# infra/ansible/templates/campustour-backend.service.j2
[Unit]
Description=CampusTour backend (staging)
After=network.target

[Service]
Type=simple
User=campustour
WorkingDirectory={{ app_dir }}/backend
# Bypass `npm start` so we don't depend on ./backend/.env existing (repo uses `node --env-file=.env`).
ExecStart=/usr/bin/node --env-file={{ env_dir }}/backend.env index.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

